# templates/default/agent_init.erb
### BEGIN INIT INFO
# Provides:          twemproxy_sentinel
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start TwemProxy Sentinel at boot time
# Description:       Starts the nodejs service to monitor redis-sentinal and update twemproxy config
### END INIT INFO

ROOT=<%= @path %>
EXEC="$ROOT/bin/redis_twemproxy_agent"
USER='root'

PID=`ps -aefw | grep "$EXEC" | grep -v " grep " | awk '{print $2}'`
USER=root

case "$1" in
    start)
        if [ ! -z $PID ]
        then
                echo "TwemProxy Redis Sentinel Agent is already running!"
        else
                echo "TwemProxy Redis Sentinel Agent Starting"
                su -c "forever start -l /dev/null -a -o /dev/null -e /dev/null -w --watchDirectory $ROOT $EXEC" $USER
                echo "TwemProxy Redis Sentinel Agent Started"
        fi
        ;;
    stop)
        if [ -z $PID ]
        then
                echo "TwemProxy Redis Sentinel Agent is not running"
        else
                echo "Stopping TwemProxy Redis Sentinel Agent"
                su -c "$FOREVER stop $EXEC" $USER
                echo "TwemProxy Redis Sentinel Agent stopped."
        fi
        ;;

    restart)
        if [ -z $PID ]
        then
                /etc/init.d/twemproxy_sentinel start
        else
                /etc/init.d/twemproxy_sentinel stop
                /etc/init.d/twemproxy_sentinel start
        fi
        ;;
    *)
        echo "Please use start or stop as first argument"
        ;;
esac
